# deploy/bootstrap/crontab
*/5 * * * * flock /var/lock/update-bootstrap bash /opt/scripts/update.sh
